# docker/supervisor_technogarden_back.conf
[program:technogarden_back]
command=php -S 0:9000 -t /var/www/TechnoGarden-backend/web
numprocs=1
autostart=true
autorestart=true
