# docker/supervisor_technogarden_front.conf
[program:technogarden_front]
priority = 4
directory = /var/www/TechnoGarden
command = node server.js
numprocs = 1
autostart = true
autorestart = true
