# docker/supervisor_mysql.conf
[program:mysqld]
priority = 2
command = mysqld_safe
numprocs = 1
autostart = true
autorestart = true
